# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/local/nginx/sbin/nginx
stopsignal=QUIT

[program:php-fpm]
command=/usr/local/php/sbin/php-fpm --nodaemonize
stopsignal=QUIT
